# package/smf/opte/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='opte'>

<service name='system/illumos/opte' type='service' version='1'>
  <create_default_instance enabled='false' />
  <dependency name='network' grouping='require_all' restart_on='none'
    type='service'>
  <service_fmri value='svc:/milestone/network:default' />
  </dependency>

  <exec_method type='method' name='start'
    exec='ctrun -l child -o noorphan,regent /opt/oxide/opte/bin/opte_method_script.sh &amp;'
    timeout_seconds='0' />
  <exec_method type='method' name='stop' exec=':kill' timeout_seconds='0' />

  <property_group name='config' type='application'>
    <propval name='underlay1' type='astring' value='' />
    <propval name='underlay2' type='astring' value='' />
  </property_group>

  <property_group name='startd' type='framework'>
    <propval name='duration' type='astring' value='transient' />
  </property_group>

  <stability value='Unstable' />

  <template>
    <common_name>
      <loctext xml:lang='C'>OPTE</loctext>
    </common_name>
    <description>
      <loctext xml:lang='C'>Oxide Packet Transformation Engine</loctext>
    </description>
  </template>
</service>

</service_bundle>
